# assets/cfdns.service
[Unit]
Description=Cloudflare DDNS updater

[Service]
Type=oneshot
ExecStart={{EXE}} update
